# infra/systemd/hybrid-frontend.service
[Unit]
Description=Hybrid Cloud — frontend (Next.js standalone)
After=network-online.target hybrid-main-api.service
Wants=network-online.target

[Service]
Type=simple
# Next.js 16 standalone in a pnpm workspace preserves the monorepo
# layout under .next/standalone/, so the real server.js lives at
# web/current/web/frontend/server.js — not at web/current/server.js.
WorkingDirectory=%h/hybrid-cloud/web/current/web/frontend
EnvironmentFile=%h/hybrid-cloud/.env.production
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hybrid-frontend

[Install]
WantedBy=default.target
