# infra/caddy/caddy.service.d/override.conf
[Service]
# Run the xcaddy build that includes the route53 DNS plugin instead of the
# apt-shipped /usr/bin/caddy binary. Reset ExecStart first because systemd
# appends rather than replaces by default.
ExecStart=
ExecStart=/usr/local/bin/caddy run --environ --config /etc/caddy/Caddyfile
ExecReload=
ExecReload=/usr/local/bin/caddy reload --config /etc/caddy/Caddyfile --force

# AWS credentials for Route53 DNS-01 challenge. File is root-owned, 0600.
EnvironmentFile=/etc/caddy/caddy.env
